# tools/deploy/taean-cctv.timer
[Unit]
Description=태안 도로 CCTV 30분마다 갱신

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
