# tools/news-audio/com.taean.newsaudio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.taean.newsaudio</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Applications/taean/tools/news-audio/run-news-audio.sh</string>
  </array>
  <!-- 매일 07:00 KST (뉴스 수집 이후) -->
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Applications/taean/tools/news-audio/news-audio.log</string>
  <key>StandardErrorPath</key><string>/Applications/taean/tools/news-audio/news-audio.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
