# tools/vps/taean-gov.timer
[Unit]
Description=태안 군청 수집 6시간마다
[Timer]
OnBootSec=3min
OnUnitActiveSec=6h
Persistent=true
[Install]
WantedBy=timers.target
