# tools/vps/taean-gov.service
[Unit]
Description=Taean 군청 수집 (한국 IP)
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
EnvironmentFile=/etc/taean.env
WorkingDirectory=__REPO__
ExecStart=/usr/bin/node tools/gov/ingest-gov.mjs
Nice=10
